# files/etc/init.d/tailscale
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

PROG=/usr/sbin/tailscaled
STATE_DIR=/etc/tailscale

start_service() {
    local enabled
    local auth_key
    local accept_routes
    local accept_dns
    local advertise_routes
    local advertise_exit_node
    local ssh
    local state_dir

    config_load 'tailscale'
    config_get enabled 'tailscale' 'enabled' '0'
    config_get auth_key 'tailscale' 'auth_key' ''
    config_get accept_routes 'tailscale' 'accept_routes' '1'
    config_get accept_dns 'tailscale' 'accept_dns' '1'
    config_get advertise_routes 'tailscale' 'advertise_routes' ''
    config_get advertise_exit_node 'tailscale' 'advertise_exit_node' '0'
    config_get ssh 'tailscale' 'ssh' '0'
    config_get state_dir 'tailscale' 'state_dir' '/etc/tailscale'

    [ "$enabled" = "1" ] || return 0

    # Criar diretório de estado
    mkdir -p "$state_dir"

    # Carregar módulo tun se necessário
    [ -d /sys/module/tun ] || modprobe tun

    # Construir argumentos
    local args="--state=$state_dir/tailscale.state --socket=$state_dir/tailscaled.sock"
    [ "$accept_routes" = "1" ] && args="$args --accept-routes"
    [ "$accept_dns" = "1" ] && args="$args --accept-dns"
    [ "$ssh" = "1" ] && args="$args --ssh"

    procd_open_instance
    procd_set_param command "$PROG" $args
    procd_set_param pidfile "$state_dir/tailscaled.pid"
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance

    # Esperar o daemon iniciar
    sleep 2

    # Conectar à tailnet se auth_key for fornecido
    if [ -n "$auth_key" ]; then
        local connect_args=""
        [ -n "$advertise_routes" ] && connect_args="$connect_args --advertise-routes=$advertise_routes"
        [ "$advertise_exit_node" = "1" ] && connect_args="$connect_args --advertise-exit-node"
        
        /usr/sbin/tailscale up --authkey="$auth_key" $connect_args
    fi
}

stop_service() {
    local state_dir
    config_load 'tailscale'
    config_get state_dir 'tailscale' 'state_dir' '/etc/tailscale'
    
    # Desconectar
    /usr/sbin/tailscale down 2>/dev/null || true
    
    # Parar daemon
    procd_kill tailscale
}

reload_service() {
    stop
    start
}
